# admin/gregor.socket
[Unit]
Description=GREGoR Dashboard gunicorn socket

[Socket]
ListenStream=/var/run/gregor.sock

[Install]
WantedBy=sockets.target
